# tpm2-rotate-swapkey.service.in
[Unit]
Description=Generate new TPM2 key for swap partition

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=$PREFIX/bin/tpm2-rotate-swapkey

PrivateTmp=true
RootEphemeral=true

[Install]
WantedBy=swaps.target
